# check_nightly_reboot/nightly-reboot.timer
[Unit]
Description=Reboot the machine at 04:00

[Timer]
OnCalendar=04:00:00
Unit=reboot.target
